# slatwall-quickstart/supervisord.conf
[supervisord]
nodaemon=true

[program:webserver]
command=catalina.sh run

[program:database]
command=mysqld --init-file=/mysql_init.sql
autorestart=false
exitcodes=0,1,2
